# etc/cron.daily/backupredis
#!/bin/bash

cp /var/lib/redis/dump.rdb ~/dump.rdb.$(date +%u)
